# projects/Eminent/filesystem/usr/lib/systemd/system/udpate-uboot.service
[Unit]
Description=update uboot

[Service]
Type=simple
ExecStart=/bin/sh -c "/usr/share/bootloader/new_uboot.sh"

[Install]
WantedBy=basic.target
